# systemd/.config/systemd/user/kmonad-keyboard.path
[Unit]
Description=Provjera da li je prikljuÄena vanjska tipkovnica

[Path]
PathExists=/dev/input/by-id/usb-Metadot_-_Das_Keyboard_Das_Keyboard-event-kbd
Unit=kmonad-keyboard.service

[Install]
WantedBy=default.target
    
